# foundations/linux/lx05-systemd/10-capstone/code/myapp-logrotate.service
# =============================================================================
# myapp-logrotate.service - Log Rotation Service for MyApp
# =============================================================================
# Triggered by myapp-logrotate.timer to clean up old logs
# Type=oneshot for short-lived tasks
# =============================================================================

[Unit]
Description=MyApp Log Rotation
After=myapp.service

[Service]
Type=oneshot
User=myapp
Group=myapp

ExecStart=/opt/myapp/bin/logrotate

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=myapp-logrotate

# Timeout
TimeoutStartSec=120

# Resource limits for batch task
CPUQuota=25%
MemoryMax=128M

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/log/myapp
