# foundations/linux/lx05-systemd/03-unit-files/code/bad-examples/secrets-inline.service
# /etc/systemd/system/bad-example.service
#
# =====================================================
# ANTI-PATTERN: Secrets in Environment= directive
# DO NOT USE THIS IN PRODUCTION!
# =====================================================
#
# This file demonstrates what NOT to do with secrets

[Unit]
Description=Bad Example - Secrets Inline (DO NOT USE!)

[Service]
Type=simple
User=appuser

# =====================================================
# BAD: Anyone can see these with 'systemctl show'!
# =====================================================
#
# Run this command to see the problem:
#   systemctl show bad-example --property=Environment
#
# Output:
#   Environment=DATABASE_URL=postgres://user:PASSWORD@localhost/db API_KEY=sk-secret123
#
Environment=DATABASE_URL=postgres://user:PASSWORD_VISIBLE_TO_ALL@localhost/db
Environment=API_KEY=sk-secret123-anyone-can-see-this
Environment=AWS_SECRET_ACCESS_KEY=wJalrXUtnFEMI/K7MDENG/bPxRfiCYEXAMPLEKEY

ExecStart=/opt/app/server

[Install]
WantedBy=multi-user.target

# =====================================================
# CORRECT APPROACH: Use EnvironmentFile
# =====================================================
#
# 1. Create secrets file with restricted permissions:
#    sudo touch /etc/myapp/secrets
#    sudo chmod 0600 /etc/myapp/secrets
#    sudo vim /etc/myapp/secrets
#
# 2. File contents:
#    DATABASE_URL=postgres://user:password@localhost/db
#    API_KEY=sk-secret123
#    AWS_SECRET_ACCESS_KEY=...
#
# 3. In Unit file:
#    EnvironmentFile=/etc/myapp/secrets
#
# 4. Now 'systemctl show' won't expose your secrets!
